# templates/nfs-sync.timer.j2
[Unit]
Description=NFS Sync Timer — ${NFS_SYNC_SCHEDULE_LABEL}
Requires=nfs-sync.service

[Timer]
OnCalendar=${NFS_SYNC_ONCALENDAR}
# Ejecutar si se perdió la ventana (ej: sistema apagado)
Persistent=true
# Variación aleatoria para evitar picos de carga exactos
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
