# templates/snapraid-sync.timer.j2
[Unit]
Description=SnapRAID Sync Timer — ${SNAPRAID_SYNC_SCHEDULE_LABEL}
Requires=snapraid-sync.service

[Timer]
OnCalendar=${SNAPRAID_SYNC_ONCALENDAR}
# Ejecutar si se perdió la ventana (ej: sistema apagado)
Persistent=true
# Pequeña variación aleatoria para evitar picos de carga exactos
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
